# bootstrap/es4/service/log/run
#!/bin/sh
# run script for logging "electroserver" daemon
exec multilog t s1048576 ./main
# that's all, folks!
